# deploy/my-agent.service
# systemd *user* service — rootless, no sudo in the deploy path.
# Installed to ~/.config/systemd/user/ by scripts/bootstrap-remote.sh.
# Runs at boot via linger. Logs: journalctl --user -u my-agent -f

[Unit]
Description=personal-assistant Telegram bot
Documentation=https://github.com/kolbeyang/my-agent

[Service]
Type=simple
WorkingDirectory=%h/my-agent
EnvironmentFile=%h/my-agent/.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/pnpm run start:prod
Restart=always
RestartSec=3
SyslogIdentifier=my-agent

[Install]
WantedBy=default.target
